# modules/buddycloud/files/buddycloud-http-api.conf
description "buddycloud HTTP API server"

start on runlevel [2345]
stop on runlevel [!2345]

console log
env NODE_ENV=production

exec node /opt/buddycloud/http-api/server.js