# modules/buddycloud/files/buddycloud-server.conf
description "buddycloud server"

start on runlevel [2345]
stop on runlevel [!2345]

console log

exec /opt/buddycloud/server/bin/buddycloud-server